# testing/admin-tester.socket
[Socket]
ListenStream=/run/admin-tester.sk

[Install]
WantedBy=sockets.target